# build-linux/deb/postinst
#!/bin/sh

set -e

if [ "$1" = "configure" ]; then
    # Create pdagent user & group
    /usr/bin/getent passwd pdagent >/dev/null || \
        /usr/sbin/adduser --system --shell /bin/false --no-create-home \
            --group pdagent

    # Fix ownership of data and log dirs
    chown -R pdagent:pdagent /var/lib/pdagent /var/log/pdagent
    chmod 777 /var/lib/pdagent/outqueue

    # Compile module .py to .pyc
    if which update-python-modules >/dev/null 2>&1; then
        update-python-modules python-pdagent.public
    fi

    # Register daemon
    update-rc.d -f pdagent defaults

    # Restart daemon
    service pdagent restart
fi

exit 0
