# build-linux/deb/postinst
#!/bin/sh
#
# Copyright (c) 2013-2014, PagerDuty, Inc. <info@pagerduty.com>
# All rights reserved.
#
# Redistribution and use in source and binary forms, with or without
# modification, are permitted provided that the following conditions are met:
#
#   * Redistributions of source code must retain the above copyright
#     notice, this list of conditions and the following disclaimer.
#   * Redistributions in binary form must reproduce the above copyright
#     notice, this list of conditions and the following disclaimer in the
#     documentation and/or other materials provided with the distribution.
#   * Neither the name of the copyright holder nor the
#     names of its contributors may be used to endorse or promote products
#     derived from this software without specific prior written permission.
#
# THIS SOFTWARE IS PROVIDED BY THE COPYRIGHT HOLDERS AND CONTRIBUTORS "AS IS"
# AND ANY EXPRESS OR IMPLIED WARRANTIES, INCLUDING, BUT NOT LIMITED TO, THE
# IMPLIED WARRANTIES OF MERCHANTABILITY AND FITNESS FOR A PARTICULAR PURPOSE
# ARE DISCLAIMED. IN NO EVENT SHALL THE COPYRIGHT HOLDER OR CONTRIBUTORS BE
# LIABLE FOR ANY DIRECT, INDIRECT, INCIDENTAL, SPECIAL, EXEMPLARY, OR
# CONSEQUENTIAL DAMAGES (INCLUDING, BUT NOT LIMITED TO, PROCUREMENT OF
# SUBSTITUTE GOODS OR SERVICES; LOSS OF USE, DATA, OR PROFITS; OR BUSINESS
# INTERRUPTION) HOWEVER CAUSED AND ON ANY THEORY OF LIABILITY, WHETHER IN
# CONTRACT, STRICT LIABILITY, OR TORT (INCLUDING NEGLIGENCE OR OTHERWISE)
# ARISING IN ANY WAY OUT OF THE USE OF THIS SOFTWARE, EVEN IF ADVISED OF THE
# POSSIBILITY OF SUCH DAMAGE.
#

set -e

install_init () {
    cp /var/lib/pdagent/scripts/pdagent.init /etc/init.d/pdagent
    chmod +x /etc/init.d/pdagent
    update-rc.d -f pdagent defaults
    service pdagent start
}

install_systemd () {
    cp /var/lib/pdagent/scripts/pdagent.service /lib/systemd/system
    systemctl enable pdagent
    if [ -z $DOCKER_DEV_ENV ]; then
        # Prevent starting pdagent in Docker environment because systemd 
        # will do that after running the Docker systemd process
        systemctl start pdagent
    fi
}

if [ "$1" = "configure" ]; then
    # Create pdagent user & group
    /usr/bin/getent passwd pdagent >/dev/null || \
        /usr/sbin/adduser --system --shell /bin/false --no-create-home \
                          --group pdagent

    # Fix ownership of data and log dirs.
    # outqueue/* world permission is set to 000 so that service keys of other
    # users cannot be seen, but set to -wx for enqueue-related directories
    # to allow any user to write events.
    chown -R pdagent:pdagent /var/lib/pdagent /var/log/pdagent
    chmod 750 /var/lib/pdagent/outqueue/err /var/lib/pdagent/outqueue/suc
    chmod 753 /var/lib/pdagent/outqueue/tmp /var/lib/pdagent/outqueue/pdq

    pycompile -p pdagent

    if which systemctl >/dev/null; then
        install_systemd
    else
        install_init
    fi
fi

exit 0
